# old-vnc-container-project/build/setup/content/.config/autostart/copyq.desktop
[Desktop Entry]
Encoding=UTF-8
Version=0.9.4
Type=Application
Name=copyq
Comment=
Exec=copyq
OnlyShowIn=XFCE;
StartupNotify=false
Terminal=false
Hidden=false
